# supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor

[program:fastapi]
command=uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
directory=/home/user/app/api
user=user
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi.err.log
environment=PYTHONPATH="/home/user/app",PYTHONUNBUFFERED=1,PATH="/usr/local/bin:/usr/bin:/bin:/home/user/.local/bin"

[program:nextjs]
command=node server.js
directory=/home/user/app/frontend
user=user
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nextjs.log
stderr_logfile=/var/log/supervisor/nextjs.err.log
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err.log
priority=100

[group:services]
programs=fastapi,nextjs,nginx
priority=5